# pkgs/workr/DEBIAN/postinst
#!/usr/bin/env bash
mkdir -p /var/srkbz/workr/jobs
systemctl enable srkbz-workr
srkbz-configure
srkbz-apply
systemctl start srkbz-workr
