# pkgs/minecraft/DEBIAN/postinst
#!/usr/bin/env bash
mkdir -p /var/srkbz/minecraft
systemctl enable srkbz-minecraft
srkbz-configure
srkbz-apply
systemctl start srkbz-minecraft
