# pkgs/minecraft/DEBIAN/postinst
#!/usr/bin/env bash
srkbz-configure
srkbz-apply
systemctl enable srkbz-minecraft
